# build/12-supervisor/my-bindfs-mounts.conf
# bindfs mounts things.
# This has to wait for "my-snaps-starter" to be done, as there's where
# main command here used is initiated.
[program:my-bindfs-mounts]
command=sh -c "while [ -z `which my-bindfs-mounts` ]; do sleep 1; done; my-bindfs-mounts --load-from-file;"
priority=111
startsecs=0