# build/12-supervisor/my-bindfs-mounts.conf
# bindfs mount things.
[program:my-bindfs-mounts]
command=/usr/local/bin/my-bindfs-mounts-init
priority=111
startsecs=0