# build/12-supervisor/my-snaps-starter.conf
# Kicks every snap package passed (a handcrafted installation)
[program:my-snaps-starter]
command=/usr/local/bin/my-snaps-starter
priority=100
startsecs=0
