# build/12-supervisor/sshd.conf
# SSH service
[program:sshd]
command=/usr/local/bin/my-bindfs-mounts-init
