# build/12-supervisor/sshd.conf
# SSH service
[program:sshd]
command=/usr/local/bin/sshd-init
