# build/11-supervisor/access-custom-user.conf
# Copy ssh keys passed (if any) from Docker's secrets mechanism into the container
[program:access-custom-user]
command=/usr/local/bin/access-custom-user
priority=1
startsecs=0
